# rc.FreeBSD
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: notifyd
# REQUIRE: DAEMON
# BEFORE:  LOGIN

. /etc/rc.subr

name="notifyd"
rcvar="notifyd_enable"
command="/usr/local/sbin/${name}"

load_rc_config $name
run_rc_command "$1"
