# rc.FreeBSD
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: stated
# REQUIRE: DAEMON
# BEFORE:  LOGIN

. /etc/rc.subr

name="stated"
rcvar="stated_enable"
command="/usr/local/sbin/${name}"

load_rc_config $name
run_rc_command "$1"
